# server/rootfs/etc/services.d/nginx/finish
#!/bin/bash

# Disable the passthrough by default
rm -f /etc/nginx/stream.d/client-passthrough.conf