# server/rootfs/etc/services.d/nginx/run
#!/command/with-contenv bash

set -e

# Start nginx in the foreground
nginx -g "daemon off;"