# server/rootfs/etc/services.d/wireguard/finish
#!/bin/bash
echo "Removing Wireguard"
ip link set wg0 down
ip link delete wg0